# scripts/com.mindbase.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mindbase</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/.local/bin/mindbase</string>
        <string>serve</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/mindbase.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/mindbase.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>DATABASE_URL</key>
        <string>postgresql+asyncpg://YOUR_USERNAME@localhost:5432/mindbase</string>
        <key>OLLAMA_URL</key>
        <string>http://localhost:11434</string>
        <key>API_PORT</key>
        <string>18002</string>
        <key>DATA_DIR</key>
        <string>/Users/YOUR_USERNAME/Library/Application Support/mindbase</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/.local/share/mindbase</string>
</dict>
</plist>
